# build-armbian/armbian-files/different-files/easepi-r2/rootfs/etc/systemd/system/eth-order.service
[Unit]
Description=EasePi Ethernet Interface Order
# Fallback mechanism: 71-net.rules (udev) is the authoritative source for
# interface naming. This service runs only if /etc/eth_order or device-tree
# eth_order exists, providing a secondary rename path.
# NOTE: systemd-udev-settle.service is intentionally NOT used here —
# it is deprecated and causes unnecessary boot delays. The eth-order.sh
# script handles missing interfaces gracefully (2>/dev/null).

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/eth-order.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
