# build-armbian/armbian-files/different-files/easepi-r2/rootfs/etc/systemd/system/sort-sync-access.service
[Unit]
Description=On-demand sort cache trigger via nginx access log monitoring
After=network.target trim_nginx.service
Wants=trim_nginx.service

[Service]
Type=simple
ExecStart=/usr/local/bin/sort_sync_access.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target