# build-armbian/armbian-files/different-files/easepi-r2/rootfs/rootfs/etc/systemd/system/set-default-gateway.service
[Unit]
Description=Set Ethernet interfaces as default gateway with priority
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set_default_gateway.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target